# os/kbd.rules
SUBSYSTEM=="usb", ATTRS{product}=="HHKB Professional", RUN+="/etc/udev/rules.d/kbd.sh"